# contrib/drfw.service
[Unit]
Description=DRFW Firewall Rules
Documentation=https://github.com/ledati16/drfw
After=network-pre.target
Before=network.target
Wants=network-pre.target

[Service]
Type=oneshot
# Apply a profile named "boot" with no confirmation prompt
# The profile must exist in /root/.local/share/drfw/profiles/boot.json
# Alternative: Use nftables.service with DRFW's "Save to System" export
ExecStart=/usr/local/bin/drfw apply boot --no-confirm
ExecReload=/usr/local/bin/drfw apply boot --no-confirm
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
